# board/common/rootfs/etc/profile.d/set_env.sh
#!/bin/sh
export PATH=$PATH:/usr/local/bin
export LANG=en_US.UTF-8
export TERM=linux